# cp2root/etc/systemd/system/rotlog.service
[Unit]
Description=rotlog

[Service]
Type=oneshot
ExecStart=/usr/bin/rotate_datalog.sh 

